# systemd/catallenya-sev.zpool.scrub.timer
[Unit]
Description=Weekly ZFS Pool Scrub Timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target