# systemd/catallenya-sev.zpool.scrub.service
[Unit]
Description=Weekly ZFS Pool Scrub
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=zpool scrub zpool